# systemd/pizero-cpuminer.service
[Unit]
Description=pizero-lottery-miner - cpuminer (SHA256d solo lottery)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust User and WorkingDirectory to match your install.
User=pi
WorkingDirectory=/home/pi/pizero-lottery-miner
ExecStart=/usr/bin/env bash scripts/run-cpuminer.sh
Restart=always
RestartSec=10
# Be a polite background citizen on the single core (let the UI stay responsive).
Nice=10

[Install]
WantedBy=multi-user.target
