# systemd/pizero-dashboard.service
[Unit]
Description=pizero-lottery-miner - LCD dashboard
After=pizero-cpuminer.service
Wants=pizero-cpuminer.service

[Service]
Type=simple
# Needs access to SPI + GPIO. Either run as root, or run as a user that is a
# member of the spi and gpio groups (default on Raspberry Pi OS: the 'pi' user).
User=pi
SupplementaryGroups=spi gpio
WorkingDirectory=/home/pi/pizero-lottery-miner
ExecStart=/usr/bin/python3 dashboard.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
